# dns-dhcp/rootfs/etc/services.d/dhcpd/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC1008
# ==============================================================================
# Start dhcpd
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================


## Run your program
exec /usr/sbin/kea-dhcp4 -c /config/kea.conf
exec /usr/sbin/kea-dhcp-ddns -c /config/kea-ddns.conf
